# install/src/main/deb/DEBIAN/prerm
#!/bin/bash

systemctl stop 'soffid-iamconsole.service' || true

systemctl disable 'soffid-iamconsole.service' || true
